# deploy/ubuntu/systemd/glucotracker-web.service
[Unit]
Description=Glucotracker FastAPI backend
After=network.target postgresql.service
Wants=postgresql.service
RequiresMountsFor=/media/megusto/storage/glucotracker

[Service]
Type=simple
User=glucotracker
Group=glucotracker
WorkingDirectory=/media/megusto/storage/glucotracker/project/backend
Environment="PATH=/media/megusto/storage/glucotracker/project/venv/bin"
EnvironmentFile=/media/megusto/storage/glucotracker/config/env
ExecStartPre=/media/megusto/storage/glucotracker/project/venv/bin/python -m alembic upgrade heads
ExecStart=/media/megusto/storage/glucotracker/project/venv/bin/python -m uvicorn glucotracker.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=10s

MemoryMax=1G
TasksMax=200

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/media/megusto/storage/glucotracker/runtime /media/megusto/storage/glucotracker/logs

StandardOutput=journal
StandardError=journal
SyslogIdentifier=glucotracker-web

[Install]
WantedBy=multi-user.target
