# deploy/ubuntu/systemd/glucotracker-worker.service
[Unit]
Description=Glucotracker background worker
After=network.target postgresql.service
Wants=postgresql.service
RequiresMountsFor=/media/megusto/storage/glucotracker

[Service]
Type=simple
User=glucotracker
Group=glucotracker
WorkingDirectory=/media/megusto/storage/glucotracker/project/backend
Environment="PATH=/media/megusto/storage/glucotracker/project/venv/bin"
EnvironmentFile=/media/megusto/storage/glucotracker/config/env
ExecStart=/media/megusto/storage/glucotracker/project/venv/bin/python -m glucotracker.workers
Restart=always
RestartSec=10s

MemoryMax=1G
TasksMax=200

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/media/megusto/storage/glucotracker/runtime /media/megusto/storage/glucotracker/logs

StandardOutput=journal
StandardError=journal
SyslogIdentifier=glucotracker-worker

[Install]
WantedBy=multi-user.target
